# install/etc/dc.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/dcnetio/bin/dc upgrade daemon
Restart=always

[Install]
WantedBy=default.target